# enable-zswap.service
[Unit]
Description=Enable zswap

[Service]
ExecStart=/usr/local/bin/enable-zswap.sh

[Install]
WantedBy=multi-user.target
